# rh_backup.service
[Unit]
Description=RomHeat backup

[Service]
Type=simple
ExecStart=/home/romheat/rh_backup.sh
